# ukrainealarm-client-gui/ukrainealarm-client.service
Description=Ukrainealarm client

[Service]
Type=simple
ExecStart=/usr/bin/ukrainealarm-client
Restart=always

[Install]
WantedBy=default.target